# deploy/systemd/asr-worker-live.service
[Unit]
Description=ASR Worker Live
After=network.target

[Service]
User=gunnar
Group=gunnar
WorkingDirectory=/srv/asr-worker
ExecStartPre=/usr/bin/mkdir -p /srv/transcribe/data/jobs/live_worker/inbox /srv/transcribe/data/jobs/live_worker/running /srv/transcribe/data/jobs/live_worker/done /srv/transcribe/data/jobs/live_worker/error
Environment=ASR_WORKER_QUEUE_BASE=/srv/transcribe/data/jobs/live_worker
Environment=ASR_WORKER_MAX_OUTSTANDING=8
Environment=ASR_WORKER_CONSUMER_ID=worker-live@1
Environment=ASR_WORKER_OPS_PORT=28110
Environment=ASR_WORKER_STATUS_OWNER=asr-worker-live
Environment=ASR_POOL_STATUS_OWNER=asr-pool
ExecStart=/srv/asr-worker/.venv/bin/python /srv/asr-worker/worker_daemon.py
Restart=always
RestartSec=2
EnvironmentFile=-/etc/asr-worker/asr-worker.env
Nice=5
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
