# deploy/systemd/asr-worker.service
# /etc/systemd/system/asr-worker.service
[Unit]
Description=ASR Worker Upload
After=network.target

[Service]
User=gunnar
Group=gunnar
WorkingDirectory=/srv/asr-worker
ExecStartPre=/usr/bin/mkdir -p /srv/omniscripta/data/upload/jobs/worker/inbox /srv/omniscripta/data/upload/jobs/worker/running /srv/omniscripta/data/upload/jobs/worker/done /srv/omniscripta/data/upload/jobs/worker/error
Environment=ASR_WORKER_QUEUE_BASE=/srv/omniscripta/data/upload/jobs/worker
Environment=ASR_WORKER_MAX_OUTSTANDING=4
Environment=ASR_WORKER_CONSUMER_ID=worker-upload@1
Environment=ASR_WORKER_OPS_PORT=28111
Environment=ASR_WORKER_STATUS_OWNER=asr-worker-batch
Environment=ASR_POOL_STATUS_OWNER=asr-pool
ExecStart=/srv/asr-worker/.venv/bin/python -m app.main
Restart=always
RestartSec=2
EnvironmentFile=-/etc/asr-worker/asr-worker.env
Nice=10
IOSchedulingClass=idle

[Install]
WantedBy=multi-user.target
